# tools/myip@.service
[Unit]
Description=Simple myip server

[Service]
ExecStart=-/usr/local/sbin/myip
StandardInput=socket
